# openbox-live/overlay/includes.chroot/etc/xdg/autostart/calamares-desktop-icon.desktop
[Desktop Entry]
Type=Application
Name=Add Calamares Desktop Icon
Exec=/bin/true
Hidden=true
NoDisplay=true
X-GNOME-Autostart-enabled=false
StartupNotify=false
