# gambian-iso/overlay/includes.chroot/etc/skel/.config/systemd/user/gamebian.service.d/override.conf
# Ensures Gamebian starts even if an older unit referenced /usr/bin/gamebian.
# The Calamares hook gamebian-web-install also patches system units under /usr/lib.
[Service]
ExecStart=
ExecStart=/usr/bin/python3 -m gamebian_app
